# src/varlink/systemd/nispor.socket
[Unit]
Description=Network Inspector(nispor) varlink service
PartOf=nispor.service

[Socket]
ListenStream=/run/nispor/nispor.so
SocketMode=0666

[Install]
WantedBy=sockets.target
